# swap.service
[Unit]
Description=swap service

[Service]
ExecStart=/home/swap.sh

[Install]
WantedBy=multi-user.target
